# application/data/fetchers/scheduler/com.fantasyai.news-snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fantasyai.news-snapshot</string>

    <!-- Absolute Python path: launchd runs with a bare environment (no user PATH).
         The script is a package module run via `-m` (see WorkingDirectory below). -->
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
        <string>-m</string>
        <string>application.data.fetchers.news</string>
        <string>snapshot</string>
    </array>

    <!-- Repo root: `python -m` puts the cwd on sys.path, so the `application` package
         resolves without an editable install. -->
    <key>WorkingDirectory</key>
    <string>/Users/willdaniel/Documents/fantasy-ai</string>

    <!-- 05:00 in the machine's LOCAL time (America/New_York) → 5am Eastern, DST-aware.
         Offset from the 4am leaguelogs job so the two daily fetchers don't overlap.
         If the Mac is asleep at 05:00, launchd runs the job at next wake. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Don't fire on load/login; only on the daily schedule (and manual kickstart). -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Logs live OUTSIDE ~/Documents: launchd cannot open files in the
         TCC-protected Documents folder (fails the job with EX_CONFIG/78). -->
    <key>StandardOutPath</key>
    <string>/Users/willdaniel/Library/Logs/fantasy-ai/news.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/willdaniel/Library/Logs/fantasy-ai/news.err.log</string>
</dict>
</plist>
